# packaging/systemd/track-system.service
[Unit]
Description=Track System - Automated Source Code Repository Tracking
Documentation=man:track-system(1)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple

# 环境配置
EnvironmentFile=/etc/track-system/track-system.env
Environment="RUST_LOG=info,track_system=debug"
Environment="RUST_BACKTRACE=1"

# 启动命令
ExecStart=/opt/track-system/bin/track-server server

# 重启策略
Restart=on-failure
RestartSec=10

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 内存限制（可选，根据实际情况调整）
# MemoryLimit=2G

# 启用 systemd 通知（可选，如果应用支持）
# Type=notify
# NotifyAccess=main

# 日志记录
StandardOutput=journal
StandardError=journal
SyslogIdentifier=track-system

[Install]
WantedBy=multi-user.target
